# dotfiles/systemd/dimmer@.service
[Unit]
Description=Screen dimmer

[Service]
ExecStart=%h/.config/i3/bin/dimmer --delay=%i
