# dotfiles/systemd/i3-companion.service
[Unit]
Description=i3 companion
PartOf=graphical-session.target

[Service]
Type=notify
ExecStart=%h/.config/i3/bin/i3-companion
Restart=on-failure
