# dotfiles/systemd/picom.service
[Unit]
Description=Compositor for X11
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/picom
Restart=on-failure
