# dotfiles/systemd/policykit-agent.service
[Unit]
Description=PolicyKit authentication agent
PartOf=graphical-session.target

[Service]
ExecStart=/usr/libexec/xfce-polkit
Restart=on-failure
