# dotfiles/systemd/polybar-weather.timer
[Unit]
Description=Weather update
PartOf=graphical-session.target

[Timer]
OnUnitActiveSec=1h
OnClockChange=true
