# dotfiles/systemd/polybar-weather.timer
[Unit]
Description=Timer for weather update
PartOf=graphical-session.target

[Timer]
OnCalendar=hourly
Persistent=true
