# dotfiles/systemd/spotify-clean.timer
[Unit]
Description=Timer for cleaning spotify extra PA streams
After=pipewire-pulse.socket

[Timer]
OnUnitActiveSec=1h
