# dotfiles/systemd/spotify-clean.timer
[Unit]
Description=Clean spotify extra PulseAudio streams
After=pipewire-pulse.socket

[Timer]
OnUnitActiveSec=1h
