# dotfiles/systemd/tray.target
[Unit]
Description=Home Manager System Tray
Requires=graphical-session-pre.target
