# dotfiles/systemd/wallpaper.timer
[Unit]
Description=Timer for wallpaper rotation
PartOf=graphical-session.target

[Timer]
OnUnitActiveSec=2h
