# dotfiles/systemd/weather.timer
[Unit]
Description=Timer for weather update
PartOf=graphical-session.target

[Timer]
OnCalendar=hourly
OnStartupSec=30
Persistent=true
