# dotfiles/systemd/dunst.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.nix-profile/bin/dunst
