# dotfiles/systemd/dunst.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/dunst
